# deploy/sv.ini

[program:microblog_web]
environment=PYTHON_EGG_CACHE=/tmp/.python-eggs/,MODE=CN
directory=/data/release/microblog
command=/usr/local/bin/gunicorn -c /deploy/gun_config.py run:app
user=user_00
autorestart=true
redirect_stderr=true
